# recipes-ai/nclawzero-demo-gemma/files/llama-server-gemma.service
[Unit]
Description=llama.cpp HTTP /v1 server for Gemma (zeroclaw local provider)
After=nclawzero-demo-gemma-fetch.service network-online.target
Wants=nclawzero-demo-gemma-fetch.service
ConditionPathExists=/var/lib/models/gemma.gguf

[Service]
Type=exec
EnvironmentFile=-/etc/nclawzero/gemma-model.env
# -c 4096 ctx; --n-gpu-layers 99 lets CUDA offload as many layers as VRAM allows
ExecStart=/usr/bin/llama-server \
    --model /var/lib/models/gemma.gguf \
    --host 127.0.0.1 \
    --port 8080 \
    -c 4096 \
    --n-gpu-layers 99
Restart=on-failure
RestartSec=5s
# Memory hardening — the demo-critical zeroclaw 5MB RAM claim lives or dies
# by keeping this process disciplined. Cap at 6GB (leaves ~2GB headroom on
# an 8GB Orin Nano Super).
MemoryMax=6G
CPUQuota=95%

[Install]
WantedBy=multi-user.target
