# recipes-ai/nclawzero-demo-gemma/files/nclawzero-demo-gemma-fetch.service
[Unit]
Description=Fetch Gemma-4-E4B GGUF for on-device inference demo
After=network-online.target srv-nclaw.mount
Wants=network-online.target
RequiresMountsFor=/srv/nclaw
ConditionPathExists=!/srv/nclaw/models/gemma-4-E4B-it-Q4_K_M.gguf

[Service]
Type=oneshot
EnvironmentFile=-/etc/nclawzero/gemma-model.env
ExecStart=/usr/libexec/nclawzero/fetch-gemma.sh
TimeoutStartSec=1800
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
